# os/board/rootddr/rootfs_overlay/etc/profile.d/input_wizard.sh
#!/bin/sh
# Input mapping profile - generated by S10inputwiz
# This file will be overwritten by the input wizard on first boot
export SDL_GAMECONTROLLERCONFIG=""
