# os/board/rootddr/rootfs_overlay/etc/profile.d/input_wizard.sh
#!/bin/sh
# Input mapping profile - sourced at login
# Populated by the input-mapping wizard on first boot
export SDL_GAMECONTROLLERCONFIG=""
